# docker/app/files/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
user=www-data
command=bash -c '/usr/bin/php /data/backend/web/bot.php && sleep 30'
autostart=true
autorestart = true