# docker/app/files/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
user=www-data
command=bash -c '/usr/bin/php /data/yii bot/run && sleep 1800'
autostart=true
autorestart = true