# docker/app/files/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
user=www-data
command=/usr/bin/php /data/yii bot/run
autostart=true
autorestart = true